# deploy/systemd/homed-recorder.service
[Unit]
Description=HOMEd Recorder Service

[Service]
ExecStart=/usr/bin/homed-recorder
Environment=LD_LIBRARY_PATH=/usr/lib
Environment=QT_PLUGIN_PATH=/usr/lib/qt5/plugins
KillSignal=SIGTERM
Restart=on-failure

[Install]
WantedBy=multi-user.target
